# services/beardog.service
[Unit]
Description=BearDog BTSP Tunnel Service
Documentation=https://github.com/ecoPrimals/beardog
After=network.target
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/beardog serve
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security
NoNewPrivileges=true
PrivateTmp=true

# Environment
Environment="RUST_LOG=info"
Environment="BEARDOG_DISCOVERY=mdns"

[Install]
WantedBy=multi-user.target

