# crates/biomeos-system/resources/biomeos.service
[Unit]
Description=biomeOS - Sovereignty-First Operating System
Documentation=https://github.com/biomeOS/biomeOS
After=network.target
Wants=network.target

[Service]
Type=notify
ExecStart=/usr/local/bin/biomeos-systemd --config /etc/biomeos/system.yaml
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StartLimitInterval=0

# Security settings
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=false
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/lib/biomeos /var/log/biomeos /run/biomeos

# Resource limits
LimitNOFILE=65536
LimitNPROC=32768

# Systemd integration
NotifyAccess=main
WatchdogSec=60
TimeoutStartSec=300
TimeoutStopSec=60

# Environment
Environment=RUST_LOG=info
Environment=BIOMEOS_MODE=systemd
EnvironmentFile=-/etc/biomeos/environment

# User and group
User=biomeos
Group=biomeos
DynamicUser=false

[Install]
WantedBy=multi-user.target 