# biomeos-root/etc/systemd/system/loamspine.service
[Unit]
Description=Loamspine Coordination Service
Documentation=https://github.com/ecoPrimals/loamspine
After=network.target beardog.service songbird.service
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/loamspine serve
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security
NoNewPrivileges=true
PrivateTmp=true

# Environment
Environment="RUST_LOG=info"
Environment="LOAMSPINE_DISCOVERY=runtime"

[Install]
WantedBy=multi-user.target

