# biomeos-root/etc/systemd/system/songbird.service
[Unit]
Description=Songbird BirdSong Discovery Service
Documentation=https://github.com/ecoPrimals/songbird
After=network.target
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/songbird daemon
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security
NoNewPrivileges=true
PrivateTmp=true

# Environment
Environment="RUST_LOG=info"
Environment="BIRDSONG_DISCOVERY=mdns"

[Install]
WantedBy=multi-user.target

