# biomeos-root/etc/systemd/system/toadstool.service
[Unit]
Description=Toadstool Data Mesh Service
Documentation=https://github.com/ecoPrimals/toadstool
After=network.target
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/toadstool daemon
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security
NoNewPrivileges=true
PrivateTmp=true

# Environment
Environment="RUST_LOG=info"
Environment="TOADSTOOL_DISCOVERY=mdns"

[Install]
WantedBy=multi-user.target

